# systemd/platform-arista.target
[Unit]
Description=Arista Platform

[Install]
Alias=sonic-platform-arista.target
WantedBy=multi-user.target
